# ansible-hyphae/roles/qdrant_snapshot_backup/templates/qdrant-backup.service.j2
[Unit]
Description=Qdrant snapshot backup
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=hyphae-backup
Group=hyphae-backup
Environment="QDRANT_URL={{ qdrant_url }}"
Environment="COLLECTION={{ qdrant_collection }}"
Environment="BACKUP_ROOT={{ qdrant_backup_root }}"
Environment="QDRANT_API_KEY={{ qdrant_api_key }}"
Environment="GCS_BUCKET={{ qdrant_gcs_bucket }}"
Environment="GCS_PREFIX={{ qdrant_gcs_prefix }}"
Environment="GOOGLE_APPLICATION_CREDENTIALS=/etc/gcloud/wif/hyphae-wif-cred.json"
Environment="GOOGLE_API_CERTIFICATE_CONFIG=/etc/gcloud/wif/certificate_config.json"
ExecStart=/usr/local/bin/qdrant_backup.sh

[Install]
WantedBy=multi-user.target
