# ansible-hyphae/roles/qdrant_snapshot_backup/templates/qdrant-backup.timer.j2
[Unit]
Description=Qdrant snapshot backup timer

[Timer]
OnCalendar=*-*-* {{ qdrant_backup_hour }}:{{ qdrant_backup_minute }}:00
Persistent=true
Unit=qdrant-backup.service

[Install]
WantedBy=timers.target
