# ansible-hyphae/roles/raspi_signage/templates/raspi-signage-watchdog.service.j2
[Unit]
Description=Watchdog for Raspberry Pi signage kiosk

[Service]
Type=oneshot
ExecStart=/bin/sh -c 'pgrep -u {{ raspi_signage_user }} -f "{{ raspi_signage_chromium_process_match }}" >/dev/null || /bin/systemctl restart raspi-signage.service'
