# ansible-hyphae/roles/raspi_signage/templates/raspi-signage-reboot.service.j2
[Unit]
Description=Nightly reboot for Raspberry Pi signage

[Service]
Type=oneshot
ExecStart=/sbin/shutdown -r now "Scheduled reboot"
