# ansible-hyphae/roles/raspi_signage/templates/raspi-signage-reboot.timer.j2
[Unit]
Description=Nightly reboot schedule for Raspberry Pi signage

[Timer]
OnCalendar=*-*-* {{ raspi_signage_reboot_time }}
Persistent=true

[Install]
WantedBy=timers.target
